# launchd/com.kernel.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.kernel.heartbeat</string>
  <key>ProgramArguments</key>
  <array>
    <string>__KERNEL_DIR__/launchd/kernel-launch.sh</string>
    <string>--heartbeat</string>
  </array>
  <key>StartCalendarInterval</key>
  <array><dict><key>Minute</key><integer>0</integer></dict></array>
  <key>StandardOutPath</key><string>__KERNEL_DIR__/kernel-memory/logs/heartbeat.out.log</string>
  <key>StandardErrorPath</key><string>__KERNEL_DIR__/kernel-memory/logs/heartbeat.err.log</string>
</dict>
</plist>
