# launchd/com.kernel.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  KERNEL nightly cleanup job (MEM-07 / MAINT-03).

  StartCalendarInterval wakes a SHORT-LIVED `node dist/index.js --cleanup` that prunes stale
  working-memory entries (quarantine/ + reflections/) and old daily logs past the retention
  window. It NEVER touches IDENTITY.md, knowledge/, finance/, or working-memory/current.md. It
  runs and exits — like --heartbeat. If the Mac was asleep, launchd runs it at the next wake.

  Schedule: nightly at 03:30 (after consolidation 03:00, before backup 04:00).

  BEFORE LOADING: replace the node path + the two /ABSOLUTE/PATH/... values (see
  launchd/README.md). node is NOT on launchd's minimal PATH — the absolute path is required.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kernel.cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/node</string>
    <string>/ABSOLUTE/PATH/TO/KernelAi/daemon/dist/index.js</string>
    <string>--cleanup</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
    <key>KERNEL_MEMORY_DIR</key>
    <string>/ABSOLUTE/PATH/TO/KernelAi/kernel-memory</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/KernelAi/kernel-memory/logs/cleanup.out.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/KernelAi/kernel-memory/logs/cleanup.err.log</string>
</dict>
</plist>
